# deploy/pocketbase.service
[Unit]
Description=PocketBase — Algorhythm auth + database
After=network.target

[Service]
Type=simple
User=root
ExecStart=/opt/pocketbase/pocketbase serve --http=127.0.0.1:8090 --dir=/opt/pocketbase/pb_data
Restart=on-failure
RestartSec=5
StandardOutput=append:/var/log/pocketbase.log
StandardError=append:/var/log/pocketbase.log

[Install]
WantedBy=multi-user.target
